# deploy/systemd/monitor.service
[Unit]
Description=AI Trading Bot Monitor Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=REPLACE_USER
WorkingDirectory=REPLACE_WORKDIR
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=REPLACE_WORKDIR/.env
ExecStart=REPLACE_WORKDIR/.venv/bin/gunicorn -c REPLACE_WORKDIR/gunicorn.conf.py monitor_app:APP
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
