# deploy/systemd/tradingbot.service
[Unit]
Description=AI Trading Bot (Paper/Live)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=REPLACE_USER
WorkingDirectory=REPLACE_WORKDIR
Environment=PYTHONUNBUFFERED=1
Environment=CUDA_VISIBLE_DEVICES=
EnvironmentFile=REPLACE_WORKDIR/.env
ExecStart=REPLACE_WORKDIR/.venv/bin/python REPLACE_WORKDIR/tradingbot.py --mode live
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
